# ops/warden-backtests.timer
[Unit]
Description=Refresh WARDENCLAW xStock backtests hourly

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Persistent=true
RandomizedDelaySec=2min
Unit=warden-backtests.service

[Install]
WantedBy=timers.target

